# systemd/vscode-tunnel-watchdog.timer
[Unit]
Description=Run VS Code Tunnel Watchdog every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
