# systemd/vscode-tunnel-watchdog.service
[Unit]
Description=VS Code Tunnel Watchdog

[Service]
Type=oneshot
ExecStart=%h/.vscode-tunnel/watchdog-tunnel.sh
